# deploy/systemd/smoki-backup.timer
[Unit]
Description=SMOKI daily local SQLite backup at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=smoki-backup.service

[Install]
WantedBy=timers.target
