# deploy/systemd/smoki-backup-full-offsite.timer
[Unit]
Description=Daily SMOKI full-project offsite backup at 04:30 MSK

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=2min
Unit=smoki-backup-full-offsite.service

[Install]
WantedBy=timers.target
