# deploy/systemd/smoki-backup-summary.timer
[Unit]
Description=SMOKI daily backup summary at 23:45 MSK

[Timer]
OnCalendar=*-*-* 23:45:00
Persistent=true

[Install]
WantedBy=timers.target
